# conf.d/brav.conf
[program:brav]
command=/gitrepo/brav.sh
autorestart=true
